# worker/com.nynm.client-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nynm.client-worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/MarshallHuff/.nvm/versions/node/v24.13.1/bin/node</string>
    <string>/Users/MarshallHuff/New General/client-hub/worker/poller.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/MarshallHuff/New General/client-hub</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/MarshallHuff/.nvm/versions/node/v24.13.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/MarshallHuff</string>
  </dict>
  <!-- Poll cadence (seconds). THE knob for reaction latency: every human interaction
       (submit, Send to Claude, approve) waits up to one full interval. A tick with no
       work is a single ~1s GET — far under Apps Script quota — so 90s is cheap.
       To change: edit here, re-copy to ~/Library/LaunchAgents, bootout + bootstrap. -->
  <key>StartInterval</key>
  <integer>90</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/MarshallHuff/New General/client-hub/worker/logs/out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/MarshallHuff/New General/client-hub/worker/logs/err.log</string>
</dict>
</plist>
